# services/homelab-bootstrap/image-cleanup.service
[Unit]
Description=Clean up old golden images

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'find /var/lib/homelab-bootstrap/images \( -name "*.img" -o -name "*.img.gz" -o -name "*.img.zst" \) -mmin +1440 -delete 2>/dev/null; echo "Image cleanup done"'
